# deploy/com.remotyy.signal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.remotyy.signal</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/remotyy</string>
        <string>signal</string>
        <string>--port</string>
        <string>9000</string>
        <string>--dev</string>
        <string>--web-dir</string>
        <string>/Users/x13/projects/remotyy/web/dist</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/remotyy-signal.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/remotyy-signal.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/x13/projects/remotyy</string>
</dict>
</plist>
